# apps/DevOpsConsole/deploy/dev-coordinator.service
[Unit]
Description=DevCoordinator loopback control API
Documentation=file:///home/DevCoordinator/skills/codex-dev-coordinator/README.md
Requires=devcoordinator-broker.service
After=network-online.target devcoordinator-broker.service
Wants=network-online.target

[Service]
Type=simple
User=holyglory
Group=holyglory
WorkingDirectory=/home/DevCoordinator
Environment=DEVCOORDINATOR_AUTHORITY=system
Environment=CODEX_AGENT_COORDINATOR_HOME=/var/lib/devcoordinator-clients/1000
ExecStart=/usr/bin/python3 /home/DevCoordinator/skills/codex-dev-coordinator/scripts/dev_coordinator.py api serve --host 127.0.0.1 --port 29876 --token-file /home/holyglory/.codex/agent-coordinator/api-token
ExecStartPost=/usr/bin/python3 /home/DevCoordinator/scripts/check_coordinator_auth_boundary.py --token-file /home/holyglory/.codex/agent-coordinator/api-token --host 127.0.0.1 --port 29876 --wait-seconds 10 --poll-interval-seconds 0.1
TimeoutStartSec=20
Restart=on-failure
RestartSec=3
# Linux gates /proc/<pid>/fd and cwd access by capability parity even for
# same-UID processes. The Console carries this one capability to bind 80/443.
# API startup clears ambient/inheritable sets before any managed child execs.
# Do not narrow CapabilityBoundingSet: this generic process manager must retain
# the manager's default ceiling for executables with legitimate file caps.
AmbientCapabilities=CAP_NET_BIND_SERVICE
# Managed dev servers are independently attributed lifecycles. Although they
# start new sessions, systemd still places descendants in this service cgroup;
# control-group shutdown would kill them when only the API is restarted.
KillMode=process

[Install]
WantedBy=multi-user.target
