# apps/DevOpsConsole/deploy/devops-console.service
[Unit]
Description=DevOps Console - TLS edge, subdomain proxy, and control panel for vr.ae
Documentation=file:///home/DevCoordinator/apps/DevOpsConsole/README.md
Requires=dev-coordinator.service
After=network-online.target dev-coordinator.service
Wants=network-online.target

[Service]
Type=simple
User=holyglory
Group=holyglory
WorkingDirectory=/home/DevCoordinator/apps/DevOpsConsole
EnvironmentFile=/home/holyglory/.config/devops-console/console.env
# EnvironmentFile values are intentionally overridden in ExecStart for every
# production-critical control/path; /usr/bin/env assignments are applied to the
# child after systemd loads the file.
ExecStartPre=/usr/bin/python3 /home/DevCoordinator/scripts/check_production_layout.py --repo-root /home/DevCoordinator --home /home/holyglory --env-file /home/holyglory/.config/devops-console/console.env --state-dir /home/holyglory/.local/state/devops-console --acme-webroot /home/holyglory/.local/state/devops-console/acme --coordinator-home /home/holyglory/.codex/agent-coordinator --token-file /home/holyglory/.codex/agent-coordinator/api-token --require-token --wait-token-seconds 10
ExecStart=/usr/bin/env DEVCOORDINATOR_ROOT=/home/DevCoordinator DEVCOORDINATOR_AUTHORITY=system COORDINATOR_AUTOSTART=0 COORDINATOR_REGISTRATION_REQUIRED=1 COORDINATOR_URL=http://127.0.0.1:29876 COORDINATOR_SCRIPT=/home/DevCoordinator/skills/codex-dev-coordinator/scripts/dev_coordinator.py COORDINATOR_TOKEN_FILE=/home/holyglory/.codex/agent-coordinator/api-token CODEX_AGENT_COORDINATOR_HOME=/var/lib/devcoordinator-clients/1000 STATE_DIR=/home/holyglory/.local/state/devops-console ACME_WEBROOT=/home/holyglory/.local/state/devops-console/acme /usr/bin/node bin/devops-console.mjs --env-file /home/holyglory/.config/devops-console/console.env
ExecStartPost=/usr/bin/python3 /home/DevCoordinator/scripts/check_console_registration_ready.py --unit devops-console.service --main-pid $MAINPID --token-file /home/holyglory/.codex/agent-coordinator/api-token --project /home/DevCoordinator --name devops-console --port 443 --host 127.0.0.1 --coordinator-port 29876 --expected-executable /usr/bin/node --expected-script bin/devops-console.mjs --env-file /home/holyglory/.config/devops-console/console.env --expected-working-directory /home/DevCoordinator/apps/DevOpsConsole --wait-seconds 80 --poll-interval-seconds 0.1
TimeoutStartSec=90
# SIGHUP reloads the TLS certificate from disk (after certbot renewal).
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3
KillMode=control-group
UMask=0077

# Bind 80/443 without root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# The console never starts the coordinator. Runtime state and credentials are
# external to the Git checkout; TLS material remains read-only.
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/holyglory/.local/state/devops-console
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
